# publisher/mqtt-trollibox@.service
[Unit]
Description=Trollibox to MQTT for player %I
After=network.target

[Service]
Type=simple
Restart=on-failure
RestartSec=2s
ExecStart=/bin/bash /opt/mqtt/mqtt-trollibox.sh %I
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target
